# systemd/gramle.service
[Unit]
Description=Gramle — Instagram date-guessing game
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
ExecStart=__NODE_BIN__ __APP_DIR__/server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
